# runit/openresty/run
#!/bin/ksh -e
exec singularity exec -w /opt/sandbox/openresty /usr/local/openresty/nginx/sbin/nginx -g "daemon off;"
