# runit/mango/run
#!/bin/ksh -e
exec 2>&1
exec singularity exec --writable /opt/sandbox/mango python3 /opt/mango/server.py --config=/etc/mango.conf
